# docker/php/laravel-fpm.conf
[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
redirect_stderr=true
